# curator1-5037961064.conf
[program:curator1-5037961064]
command=python3 /home/curator1/curator1-5037961064/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator1/curator1-5037961064/test.err.log
stdout_logfile=/home/curator1/curator1-5037961064/test.out.log

